# systemd/crypto-models-buckets.timer
[Unit]
Description=Nightly bucket update

[Timer]
OnCalendar=*-*-* 01:50:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
