# systemd/crypto-models-train.timer
[Unit]
Description=Daily crypto-models training

[Timer]
OnCalendar=*-*-* 02:30:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
